# crontab

# m h  dom mon dow   command
@reboot /var/www/restartpm2.sh
*/1 * * * * bash /root/test-restart-server.sh
